# daily-wallpaper-arch/dynamic-wallpaper.timer
[Unit]
Description=Run Dynamic Wallpaper Script Hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target