# daily-wallpaper-arch/dynamic-wallpaper.timer
[Unit]
Description=Run Dynamic Wallpaper Script Hourly

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target