# daily-wallpaper-arch/dynamic-wallpaper.timer
[Unit]
Description=Run Dynamic Wallpaper Script Hourly

[Timer]
OnBootSec=1min
OnCalendar=*-*-* *:22:00

[Install]
WantedBy=timers.target