# daily-wallpaper-arch/setup.install
post_install() {
    systemctl enable --now dynamic-wallpaper.timer
}

post_upgrade() {
    systemctl enable --now dynamic-wallpaper.timer
}

pre_upgrade() {
    systemctl disable --now dynamic-wallpaper.timer
}

post_remove() {
    systemctl disable --now dynamic-wallpaper.timer
}